# xmonad/autostart/xcompmgr.desktop
[Desktop Entry]
Name=xcompmgr
Comment=Composite for X
Exec=xcompmgr
Terminal=false
Type=Application
StartupNotify=false
